# docker/service/node/run
#!/bin/sh -eu
cd /data
exec oasis-node --config "$OASIS_NODE_CONFIG"
