# docker/service/gateway/run
#!/bin/sh -eu
sv start node
cd /data
exec chpst -u www-data oasis-rosetta-gateway
