# docker/service/gateway/run
#!/bin/sh -eu
sv start node
cd /data
exec oasis-rosetta-gateway
